# vicsensor.service
[Unit]
Description=Vic-sensord

[Service]
Type=simple
ExecStart=/usr/bin/vic-sensor
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target